# jobs/cron.monthly/run_monthly
#!/bin/sh

mkdir -p /data/logs
DATE=`date +%Y-%m-%d:%H:%M:%S`
echo $DATE >> /data/logs/monthly.log
nohup /usr/bin/python3 /data/stock/jobs/month_job.py > /dev/null 2>/data/logs/month.0.log &
